# coresight-agent/coresight-agent.service
[Unit]
Description=CoreSight Monitoring Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/coresight/coresight-agent
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/root/coresight/coresight-agent/.env
ExecStart=/root/coresight/coresight-agent/venv/bin/python3 agent.py
Restart=always
RestartSec=10
StandardOutput=append:/root/coresight/coresight-agent/logs/agent.log
StandardError=append:/root/coresight/coresight-agent/logs/agent.log
LogsDirectory=/root/coresight/coresight-agent/logs

[Install]
WantedBy=multi-user.target 